# profile.d/20directory.sh
##
# directory.sh
##

##
# Directory helpers
##

# navigation
alias --  -='cd -'
alias    ..='cd ..'
alias    up='cd ..'
function pd { pushd ${1:-~} ;  }

# creation
alias    md='mkdir'
alias    rd='rmdir'
alias   rdf='rm -rf'

##
# EOF
##

