# profile.d/school.sh
##
# school.sh
##

function mark { ssh mark"$@"@localhost:2222 ; }

##
# EOF
##
